# systemd/ais-server.service
[Unit]
Description=JLBMaritime AIS-Server
Documentation=https://github.com/JLBMaritime/AIS-Server
After=network-online.target tailscaled.service NetworkManager.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
WorkingDirectory=/opt/ais-server
ExecStart=/opt/ais-server/.venv/bin/python -m ais_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=2
# Allow up to 60s for a watchdog ping (the netcheck/recovery cycle can take
# ~20s under a Wi-Fi bounce).
WatchdogSec=60
TimeoutStopSec=15

# Run as root so the wifi / systemctl helpers work; drop privs for everything
# else via Python-level sandboxing.  If you prefer, create a `ais` user and
# set User=ais + grant sudoers for nmcli/systemctl.
User=root
Group=root

# Stability / resource guard-rails.  TasksMax raised from the systemd default
# of 512 because Flask-SocketIO in async_mode="threading" can briefly spike
# the thread count under reconnect storms.
MemoryHigh=512M
MemoryMax=1G
TasksMax=2048
LimitNOFILE=65536

# Hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/etc/ais-server /var/lib/ais-server

# Logs go to stdout -> journald (single source of truth – use
# `journalctl -u ais-server -f` to follow).  Do NOT add StandardOutput=append:
# here, that conflict with Python rotation has caused 24h freezes in the past.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
